# packages/python/mopidy/system.d/mopidy.service
[Unit]
Description=Mopidy music server
After=avahi-daemon.service
After=dbus.service
After=network-online.service
After=nss-lookup.target
After=remote-fs.target
After=sound.target
Wants=network-online.target

[Service]
ExecStartPre=/bin/mkdir -p /storage/.cache/mopidy
ExecStart=python3 /usr/bin/mopidy --config /usr/share/mopidy/conf.d:/storage/.config/mopidy/mopidy.conf

[Install]
WantedBy=multi-user.target