# packages/graphics/cog/system.d/cog.service
[Unit]
Description=Cog service
Wants=weston.service
After=weston.service
After=mopidy.service

[Service]
Environment=XDG_RUNTIME_DIR=/var/run/0-runtime-dir/
Environment=WAYLAND_DISPLAY=wayland-1
ExecStart=/usr/bin/cog-launcher
Restart=on-failure
RestartSec=10

[Install]
WantedBy=graphical.target
