# packages/graphics/cog/system.d/cog.service
[Unit]
Description=Cog service
Requires=weston.service
After=weston.service
After=mopidy.service

[Service]
Environment=XDG_RUNTIME_DIR=/var/run/0-runtime-dir/
Environment=WAYLAND_DISPLAY=wayland-0
ExecStart=/usr/bin/cog-launcher

[Install]
WantedBy=graphical.target
