# packages/python/hifistreamer-app/system.d/hifistreamer-app.service
[Unit]
Description=HiFiStreamer Application
After=mopidy.service
After=sound.target

[Service]
ExecStart=python3 -m hifistreamer-app

[Install]
WantedBy=multi-user.target
