# packages/sysutils/shairmopd/system.d/shairmopd.service
[Unit]
After=shairport-sync.service
After=mopidy.service

[Service]
ExecStart=/usr/bin/shairmopd

[Install]
WantedBy=multi-user.target
